# listen-engine/listen-engine.service
[Unit]
Description=Listen Engine Service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=listen-engine
Group=listen-engine
Restart=always
RestartSec=1
Environment=IS_SYSTEMD_SERVICE=1
EnvironmentFile=/etc/listen-engine/environment
WorkingDirectory=/opt/listen-engine
ExecStart=/usr/local/bin/engine

# Configure logging to both file and journal
StandardOutput=append:/var/log/listen-engine/service.log+journal+console
StandardError=append:/var/log/listen-engine/error.log+journal+console
SyslogIdentifier=listen-engine

[Install]
WantedBy=multi-user.target
